# reminder-app.service
[Unit]
Description=Reminder App Background Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/reminder-app --background
Restart=on-failure
RestartSec=5
User=%u

[Install]
WantedBy=default.target